# docs/namingo-epp-reload.service
[Unit]
Description=Reload Namingo EPP server after cert update

[Service]
Type=oneshot
ExecStart=/bin/systemctl reload epp